# contrib/launchd/com.askdba.mysql-mcp-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.askdba.mysql-mcp-server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/mysql-mcp-server</string>
		<string>--config</string>
		<string>/usr/local/etc/mysql-mcp-server/config.yaml</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MYSQL_MCP_HTTP</key>
		<string>1</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>StandardOutPath</key>
	<string>/usr/local/var/log/mysql-mcp-server.log</string>
	<key>StandardErrorPath</key>
	<string>/usr/local/var/log/mysql-mcp-server.err</string>
</dict>
</plist>
